# ps3netsrv/rootfs/etc/s6-overlay/s6-rc.d/init-ps3netsrv/up
/etc/s6-overlay/s6-rc.d/init-ps3netsrv/run