# oscam/root/etc/services.d/oscam/run
#!/usr/bin/with-contenv bash

exec \
  /usr/bin/oscam -c /config/oscam

#exec \
#  s6-setuidgid abc /usr/bin/oscam -c /config/oscam
