# oscam/root/etc/services.d/oscam/run
#!/usr/bin/with-contenv bash
exec \
  s6-setuidgid root /usr/bin/oscam -c /config/oscam
#s6-setuidgid abc /usr/bin/oscam -c /config/oscam
