# ps3netsrv/rootfs/etc/s6-overlay/s6-rc.d/ps3netsrv/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Ps3NetSrv
# Runs the ps3netsrv
# ==============================================================================

bashio::log.info "Starting Ps3NetSrv...."
/usr/local/bin/ps3netsrv /share/ps3netsrv
